# supervisord.conf
[supervisord]
nodaemon=false

[program:node]
command=npm start
directory=/socket
user=root
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
autorestart=true